# systemd/multibot.timer
[Unit]
Description=artspambot all the time

[Timer]
OnBootSec=10min
OnCalendar=*:*:0/10

[Install]
WantedBy=timers.target
